# files/system/usr/lib/systemd/system/config-manager.service
[Unit]
Description=Config Initialization
Before=display-manager.service

[Service]
Type=oneshot
ExecStart=/usr/bin/python3 /usr/lib/config-manager/update.py
RemainAfterExit=no

[Install]
WantedBy=multi-user.target